# yt-learner-worker.service
[Unit]
Description=yt-learner worker
After=network.target

[Service]
Type=simple
User=deepanshu
EnvironmentFile=/home/deepanshu/.env
WorkingDirectory=/home/deepanshu/projects/yt-learner
ExecStart=/home/deepanshu/.local/bin/uv run yt-learner-worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
